# packages/bsp/getek/lib/systemd/system/display-manager.service
[Unit]
Description=Display Manager
After=systemd-user-sessions.service

[Service]
ExecStart=/usr/bin/X

[Install]
Alias=display-manager.service

